# deploy/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:backend]
command=/bin/sh -c "java -jar /app/backend/app.jar --spring.profiles.active=docker --server.port=8080"
priority=10
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:watchbug]
command=/bin/sh -c "PORT=5000 python /app/watchbug-app/app.py"
directory=/app/watchbug-app
priority=20
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
priority=30
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0